# telegram-bot.service
[Unit]
Description=Telegram Ollama Bot
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/telegram-bot
ExecStart=/opt/telegram-bot/venv/bin/python bot.py
Restart=on-failure
RestartSec=5

# Put secrets in /etc/telegram-bot.env (KEY=value, one per line)
EnvironmentFile=-/etc/telegram-bot.env

[Install]
WantedBy=multi-user.target
